# chat-mcp-setup/host/macos/mcp-gateway/com.mcp-gateway.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcp-gateway.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/caddy</string>
        <string>run</string>
        <string>--config</string>
        <string>$HOME/mcp-gateway/Caddyfile</string>
        <string>--adapter</string>
        <string>caddyfile</string>
    </array>

    <key>WorkingDirectory</key>
    <string>$HOME/mcp-gateway</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>$HOME</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LAN_MCP_TOKEN</key>
        <string>__TOKEN_REPLACE_ME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/mcp-gateway/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/mcp-gateway/stderr.log</string>
</dict>
</plist>
