# mcp-bundle/supervisord.conf
; Single-container fan-out: one supervisord runs 5 FastMCP processes,
; each on its own port. All bind 0.0.0.0 via MCP_HOST env so Coolify
; Traefik / sibling containers can reach them on the coolify network.
; SECURITY (audit I1): set MCP_SHARED_SECRET in the Coolify env to require an
; `Authorization: Bearer <secret>` on every request (opt-in via mcp/_mcp_auth.py;
; unset = no auth, unchanged). The SAME secret must be set on the alphalens
; service so its MCP client sends it.
;
; Logs all funnel to stdout/stderr → Coolify's container log UI shows
; the joined stream (one container, all 5 services).

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; -------- per-service programs --------
; Stdout/stderr go to /dev/fd/1 + /dev/fd/2 → docker logs visible.
; autorestart catches transient API crashes; startsecs prevents flap loop.

[program:fred]
command=python -u /app/mcp/fred_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",FRED_MCP_PORT="3035"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

[program:tradier]
command=python -u /app/mcp/tradier_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",TRADIER_MCP_PORT="3036"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:polymarket]
command=python -u /app/mcp/polymarket_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",POLYMARKET_MCP_PORT="3033"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:stock-price]
command=python -u /app/mcp/stock_price_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",STOCK_MCP_PORT="3032"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:financialjuice]
command=python -u /app/mcp/financialjuice_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",FINANCIALJUICE_MCP_PORT="3034"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sa-rss]
command=python -u /app/mcp/seekingalpha_rss_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",SA_RSS_MCP_PORT="7050"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sec-edgar]
command=python -u /app/mcp/sec_edgar_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",SEC_EDGAR_MCP_PORT="3037",SEC_EDGAR_CONTACT="user001@token-gateway.com"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sec-flows]
command=python -u /app/mcp/sec_flows_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",SEC_FLOWS_MCP_PORT="3038",SEC_EDGAR_CONTACT="user001@token-gateway.com"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; -------- P1 free-source servers (sentiment/shorts/finnhub/defeatbeta/gdelt) --------
; shorts inherits FINRA_API_CLIENT_ID/SECRET from the container env (Coolify);
; finnhub inherits FINNHUB_API_KEY; defeatbeta inherits DEFEATBETA_CACHE_DIR.
; Do NOT write any credentials in this file.

[program:sentiment]
command=python -u /app/mcp/sentiment_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",SENTIMENT_MCP_PORT="3039"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:shorts]
command=python -u /app/mcp/shorts_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",SHORTS_MCP_PORT="3040"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:finnhub]
command=python -u /app/mcp/finnhub_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",FINNHUB_MCP_PORT="3041"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:defeatbeta]
command=python -u /app/mcp/defeatbeta_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",DEFEATBETA_MCP_PORT="3042",DEFEATBETA_CACHE_DIR="/app/.defeatbeta-cache"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:gdelt]
command=python -u /app/mcp/gdelt_mcp.py
directory=/app
environment=MCP_HOST="0.0.0.0",GDELT_MCP_PORT="3043"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
